# infra/systemd/novelclaw-multiagent.service
[Unit]
Description=CoLong MultiAgent Studio
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/novelclaw/apps/multiagent
EnvironmentFile=/etc/novelclaw/multiagent.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONUTF8=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/novelclaw/.venv/bin/uvicorn local_web_portal.app.main:app --host 127.0.0.1 --port 8011 --workers 1
Restart=always
RestartSec=5
KillMode=process

[Install]
WantedBy=multi-user.target
