# infra/systemd/novelclaw-portal.service
[Unit]
Description=CoLong Auth Portal
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/novelclaw/apps/auth-portal
EnvironmentFile=/etc/novelclaw/auth-portal.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONUTF8=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/novelclaw/.venv/bin/uvicorn local_web_portal.app.main:app --host 127.0.0.1 --port 8010 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
